# package/lean/luci-app-ardnspod/root/etc/hotplug.d/iface/99-ardnspod
#!/bin/sh
# wan ifup
[ "$ACTION" = "ifup" ] && /etc/init.d/ardnspod start
