# package/lean/luci-app-ardnspod/root/etc/uci-defaults/ardnspod
#!/bin/sh

/etc/init.d/ardnspod enable && /etc/init.d/ardnspod start

rm -f /tmp/luci-indexcache
exit 0
